# systemd/opagentd.service
[Unit]
Description=opagentd - LLM-driven system agent daemon
Documentation=https://github.com/anomalyco/opagentd
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/opagentd
Restart=on-failure
RestartSec=5s
User=root
RuntimeDirectory=opagentd
RuntimeDirectoryMode=0700
LogsDirectory=opagentd
LogsDirectoryMode=0700
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
